# health-dashboard.service
[Unit]
Description=Health Dashboard Server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/keo/repos/health-dashboard
ExecStart=/usr/bin/python3 /home/keo/repos/health-dashboard/serve.py
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
